# infra/trader-daemon.service
[Unit]
Description=Trader Daemon
After=network.target

[Service]
Type=simple
User=trader
Group=trader
WorkingDirectory=/opt/trader
ExecStart=/usr/bin/npx ts-node daemon.ts start
Restart=always
RestartSec=30
Environment=NODE_ENV=production
Environment=HOME=/home/trader
Environment=PATH=/home/trader/.local/bin:/usr/local/bin:/usr/bin:/bin

# Load environment variables from .env
EnvironmentFile=/opt/trader/.env

# Logging
StandardOutput=append:/opt/trader/state/logs/daemon.log
StandardError=append:/opt/trader/state/logs/daemon.log

[Install]
WantedBy=multi-user.target
